# board/edufw/rootfs-overlay/etc/periodic/daily/edufw_blacklist_download_all
../../../sbin/edufw_blacklist_download_all